# dev-notes/systemd-examples/my-socket@.service
[Unit]
Description=Socket activated echo service

[Service]
ExecStart=/bin/bash -c "echo 'Hello from Systemd Socket! Type something and press Enter:'; cat"
StandardInput=socket
StandardOutput=socket
