# dev-notes/systemd-examples/my-service.service
[Unit]
Description=My Simple Hello World Service

[Service]
Type=oneshot
ExecStart=/home/jollof/coding/dev-setup/systemd-examples/hello_service.sh
